# com.env.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.env.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/YOUR_USERNAME/setup_dotfiles.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USERNAME/.zshrc</string>
        <string>/Users/YOUR_USERNAME/.config/starship.toml</string>
        <string>/Users/YOUR_USERNAME/Library/Application Support/com.mitchellh.ghostty/config.ghostty</string>
        <string>/Applications</string>
        <string>/opt/homebrew/Caskroom</string>
        <string>/opt/homebrew/Cellar</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
